# services/systemd/securebot-daily.service
[Unit]
Description=SecureBot Daily Report

[Service]
Type=oneshot
User=tasker0
ExecStart=/home/tasker0/securebot/services/scripts/daily_report.sh
